# files/lib/systemd/system/bluelab-first-boot.service
[Unit]
Description=BlueLab First Boot Configuration
Documentation=https://github.com/your-username/bluelab
# Start after basic system is ready, but don't require GUI services that might not exist
After=network-online.target docker.service multi-user.target
Wants=network-online.target docker.service
# Only run if first boot hasn't completed
ConditionPathExists=!/var/lib/bluelab/.first-boot-complete
# Don't run in rescue mode or during installation
ConditionKernelCommandLine=!rescue
ConditionKernelCommandLine=!inst.stage2

[Service]
Type=oneshot
ExecStart=/usr/bin/bluelab-first-boot
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console
# Give it time to complete
TimeoutStartSec=1800
# Run as root (needed for system configuration)
User=root
Group=root
# Set up environment
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=BLUELAB_AUTO_INSTALL=true
Environment=BLUELAB_USERNAME=bluelab
Environment=BLUELAB_PASSWORD=bluelab123
Environment=BLUELAB_HOSTNAME=bluelab
Environment=BLUELAB_TIMEZONE=America/New_York
# Restart on failure
Restart=on-failure
RestartSec=30
StartLimitBurst=3

[Install]
WantedBy=multi-user.target