# files/lib/systemd/system/bluelab-first-boot.service
[Unit]
Description=BlueLab First Boot Configuration
Documentation=https://github.com/your-username/bluelab
# Wait for system initial setup to complete first
After=network-online.target docker.service graphical-session.target gnome-initial-setup.service initial-setup.service
Wants=network-online.target
Requires=docker.service
# Only run if first boot hasn't completed
ConditionPathExists=!/var/lib/bluelab/.first-boot-complete
# Don't run in rescue mode or during installation
ConditionKernelCommandLine=!rescue
ConditionKernelCommandLine=!inst.stage2
# Don't run if system initial setup is still pending
ConditionPathExists=!/run/initial-setup
ConditionPathExists=!/var/lib/gnome-initial-setup/state

[Service]
Type=oneshot
ExecStart=/usr/bin/bluelab-first-boot
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console
# Give it time to complete
TimeoutStartSec=1800
# Run as root (needed for system configuration)
User=root
Group=root
# Set up environment
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=BLUELAB_AUTO_INSTALL=true
# Restart on failure
Restart=on-failure
RestartSec=30
StartLimitBurst=3

[Install]
WantedBy=graphical.target